# packages/xmedia/torrent/aceproxy/system.d/aceproxy.service
[Unit]
Description=AceProxy server
Requires=network-online.target
Wants=acestream.service
After=network-online.target acestream.service
Before=kodi.service vdr.service
ConditionPathExists=/storage/.cache/services/aceproxy.conf

[Service]
ExecStartPre=/usr/bin/aceproxy-config
ExecStart=/storage/.config/aceproxy/acehttp.py
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID
TimeoutStopSec=1
Restart=always
RestartSec=2
StartLimitInterval=0

[Install]
WantedBy=multi-user.target
